# a2im.service
[Unit]
Description=AppleMidi2IPMidiBridge

[Service]
ExecStart=/usr/bin/a2im

[Install]
WantedBy=multi-user.target
